# jobs/dea/templates/correct_uhurufs_bindings.cron
1,16,31,46 * * * * /var/vcap/jobs/dea/bin/correct_uhurufs_bindings.sh
