# jobs/dea/templates/correct_uhurufs_bindings.cron
*/15 * * * * /var/vcap/jobs/dea/bin/correct_uhurufs_bindings.sh